# misc/linux/wopded.service.in
[Unit]
Description=This service spawns a World Of Padman Internet server

[Service]
User=wop
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/@OUTPUTNAME@ +set dedicated 2 +exec server-settings.cfg +vstr nextmap
Restart=on-abort

[Install]
WantedBy=multi-user.target
